# debian/prerm
#!/bin/sh

set -e

case "$1" in
    remove|upgrade|deconfigure)
		update-rc.d -f hamonikr-system remove

        if which systemctl >/dev/null 2>&1
        then
            if [ $(systemctl is-enabled ubuntu-system-adjustments) = "disabled" ] ; then
                systemctl enable ubuntu-system-adjustments
            fi        
            systemctl daemon-reload
        fi

        # Apply HamoniKR OS infomation in system
        cp -f /etc/hamonikr/issue.bak /etc/issue
        cp -f /etc/hamonikr/issue.net.bak /etc/issue.net
        cp -f /etc/hamonikr/lsb-release.bak /etc/lsb-release
        cp -f /etc/hamonikr/os-release.bak /etc/os-release
        cp -f /etc/hamonikr/os-release.bak /usr/lib/os-release

        rm -f /etc/hamonikr/issue.bak
        rm -f /etc/hamonikr/issue.net.bak
        rm -f /etc/hamonikr/lsb-release.bak
        rm -f /etc/hamonikr/os-release.bak

        rm -f /etc/grub.d/40_hamonikr

        ;;
    failed-upgrade)
        ;;
    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

